# caom2-tapdb/pgdb-init.service
[Unit]
After=network.service

[Service]
ExecStart=/var/lib/pgsql/10/do-init.sh

[Install]
WantedBy=default.target
